# docker/osrm/service/routed/run
#!/bin/sh

exec osrm-routed --port 5000 --shared-memory --max-table-size 100000 --max-viaroute-size 1000 --algorithm $ALGORITHM 2>&1
